# com.tennisbooker.scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tennisbooker.scheduler</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/brandonvincent/Documents/personal/tennis-booker/.venv/bin/python</string>
        <string>/Users/brandonvincent/Documents/personal/tennis-booker/scheduler.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/brandonvincent/Documents/personal/tennis-booker</string>

    <!-- Run once a day at 5:50 PM (10 min before first slot) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>17</integer>
        <key>Minute</key>
        <integer>50</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/brandonvincent/Documents/personal/tennis-booker/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/brandonvincent/Documents/personal/tennis-booker/launchd_stderr.log</string>

    <!-- Restart if it crashes -->
    <key>KeepAlive</key>
    <false/>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
